# inittab
# at startup, redirect all logs to stdout for docker
::sysinit:ln -s /dev/stdout /dev/log

# AT
# E0 (echo off)
# Q1 (no results codes)
# &D2 (normal DTR operations)
# &C1 (normal CD operations)
# S0=1 (Auto Answer after 1 ring)
# \015 ???

ttyUSB0:2345:respawn:/sbin/agetty --extract-baud --timeout 300 --init-string 'ATE0Q1&D2&C1S0=1\015' -J -w ttyUSB0 9600,2400,1200,300
